# group-manager-server/conf/ukui-group-manager.service
[Unit]
Description=group-manager-server

[Service]
User=root
ExecStart=/usr/bin/group-manager-server

[Install]
WantedBy=multi-user.target
