# debian/ukui-control-center.postinst
#!/bin/sh

set -e

glib-compile-schemas /usr/share/glib-2.0/schemas/

# 启动服务
#systemctl enable ukui-group-manager.service
#systemctl start  ukui-group-manager.service

chown root:root /usr/bin/checkuserpwd
chmod u+s /usr/bin/checkuserpwd

#DEBHELPER#
